# com.zubinshourie.morning-run.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zubinshourie.morning-run</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ZYN/Desktop/site-prototype/morning-run.sh</string>
    </array>

    <!-- Run at 8:00 AM every day (after the 7:09 AM Cowork brief task finishes) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ZYN/Desktop/morning-run.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ZYN/Desktop/morning-run.log</string>

    <!-- Do NOT run immediately on login (only at 8 AM) -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
